# systemd/persistence.timer
[Unit]
Description=Persistence Timer

[Timer]
Unit=persistence.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
